# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1366x768x24
autorestart=true
user=alpine
priority=100

[program:x11vnc]
command=/usr/bin/x11vnc -xkb -noxrecord -noxfixes -noxdamage -display :1 -nopw -wait 5 -shared -permitfiletransfer -tightfilexfer
user=alpine
autorestart=true
priority=200

[program:openbox]
environment=HOME="/home/alpine",DISPLAY=":1",USER="alpine"
command=/usr/bin/openbox
user=alpine
autorestart=true
priority=300

[program:wget]
user=alpine
command=/usr/bin/wget https://osbot.org/mvc/get -O osbot.jar
autorestart=false
priority=400

[progran:osbot]
user=alpine
environment=DISPLAY=":1"
command=/usr/bin/java -jar /home/alpine/osbot.jar
priority=500